# docker/s6-rc.d/laravel-queue/run
#!/command/execlineb -P
s6-setuidgid www-data
php /var/www/html/artisan queue:work --tries=3 --max-time=3600
