# docker/s6-rc.d/laravel-scheduler/run
#!/command/execlineb -P
s6-setuidgid www-data
cd /var/www/html
php /var/www/html/artisan schedule:work
