# docker/s6-rc.d/php-fpm/run
#!/command/execlineb -P
exec /usr/local/sbin/php-fpm --nodaemonize
